# config/osquery.conf
{
  "options": {
    "logger_plugin": "filesystem",
    "logger_path": "/var/log/osquery",
    "disable_logging": "false",
    "log_result_events": "true",
    "schedule_splay_percent": "10",
    "events_expiry": "3600",
    "database_path": "/var/osquery/osquery.db",
    "verbose": "false",
    "worker_threads": "2",
    "enable_file_events": "true",
    "disable_events": "false",
    "host_identifier": "hostname"
  },

  "schedule": {
    "processes": {
      "query": "SELECT pid, name, path, cmdline, uid, gid, parent, start_time, resident_size FROM processes;",
      "interval": 60,
      "description": "All running processes — detects new/suspicious process launches"
    },
    "listening_ports": {
      "query": "SELECT lp.pid, lp.port, lp.protocol, lp.address, p.name, p.path FROM listening_ports lp JOIN processes p ON lp.pid = p.pid WHERE lp.address != '127.0.0.1';",
      "interval": 120,
      "description": "Network listeners excluding localhost — detects backdoors and C2 beacons"
    },
    "logged_in_users": {
      "query": "SELECT type, user, host, time, pid FROM logged_in_users;",
      "interval": 60,
      "description": "Active sessions — detects unauthorized access"
    },
    "open_sockets": {
      "query": "SELECT pid, remote_address, remote_port, local_address, local_port, protocol, p.name, p.path FROM process_open_sockets pos JOIN processes p ON pos.pid = p.pid WHERE remote_address != '' AND remote_address != '::' AND remote_address != '0.0.0.0';",
      "interval": 60,
      "description": "Active network connections — critical for C2 and exfil detection"
    },
    "shell_history": {
      "query": "SELECT username, command, history_file FROM shell_history WHERE command != '';",
      "interval": 600,
      "description": "Command history — detects post-exploitation activity. NOTE: 10min interval (expensive query)"
    },
    "crontab": {
      "query": "SELECT event, minute, hour, day_of_month, month, day_of_week, command, path FROM crontab;",
      "interval": 300,
      "description": "Scheduled tasks — persistence mechanism detection"
    },
    "startup_items": {
      "query": "SELECT name, path, source, status, username FROM startup_items;",
      "interval": 300,
      "description": "Login items and launch agents — persistence detection"
    },
    "launchd_entries": {
      "query": "SELECT name, label, program, program_arguments, run_at_load FROM launchd WHERE run_at_load = '1';",
      "interval": 300,
      "description": "LaunchDaemons/Agents set to auto-run — macOS persistence"
    },
    "browser_plugins": {
      "query": "SELECT * FROM browser_plugins;",
      "interval": 3600,
      "description": "Browser extensions — can be used for credential theft"
    },
    "disk_encryption": {
      "query": "SELECT * FROM disk_encryption WHERE encrypted = 1;",
      "interval": 3600,
      "description": "FileVault status — compliance check"
    },
    "sip_config": {
      "query": "SELECT * FROM sip_config;",
      "interval": 3600,
      "description": "System Integrity Protection — should always be enabled"
    },
    "user_ssh_keys": {
      "query": "SELECT uid, path, encrypted FROM user_ssh_keys;",
      "interval": 3600,
      "description": "SSH keys — monitors for unauthorized key creation"
    }
  },

  "file_paths": {
    "homes": [
      "/Users/%/Documents/%%",
      "/Users/%/Downloads/%%",
      "/Users/%/Desktop/%%",
      "/Users/%/.ssh/%%"
    ],
    "system_binaries": [
      "/usr/bin/%%",
      "/usr/local/bin/%%",
      "/opt/homebrew/bin/%%"
    ],
    "launch_agents": [
      "/Library/LaunchAgents/%%",
      "/Library/LaunchDaemons/%%",
      "/Users/%/Library/LaunchAgents/%%"
    ]
  },

  "file_accesses": ["homes", "system_binaries", "launch_agents"]
}
