# rules/sigma/macOS/launch_agent_persistence.yml
title: LaunchAgent Persistence Created
id: 03574d52-0001-5070-bb2e-5c998f480a76
status: stable
description: Detects creation of macOS LaunchAgent plist files for persistence
author: SecurityScarletAI
date: 2026/05/03
logsource:
    category: file
detection:
    selection:
        event_type: start
        file_path|contains: LaunchAgents
    condition: selection
level: high
tags:
    - attack.persistence
    - attack.ta0003
    - attack.t1547