# rules/sigma/launch_agent_persistence.yml
title: New Launch Agent Persistence
id: scarlet-002
status: experimental
description: Detects creation of LaunchAgent plist files for persistence
author: SecurityScarletAI
date: 2026/03/21
logsource:
    category: file
detection:
    selection:
        event_type: start
        file_path|contains: LaunchAgents
    condition: selection
level: high
tags:
    - attack.persistence
    - attack.t1547